# private_dot_config/powershell/profile.d/.chezmoiremove
fnm.ps1.tmpl
fnm.ps1
vfox.ps1.tmpl
vfox.ps1
