# private_dot_config/powershell/profile.d/carapace.ps1
carapace _carapace powershell | Out-String | Invoke-Expression
