# private_dot_config/powershell/profile.d/fnm.ps1.tmpl
{{ output "fnm" "env" "--use-on-cd" "--shell" "power-shell" }}
