# private_dot_config/powershell/profile.d/mise.ps1.tmpl
{{ if lookPath "mise" -}}
{{ output "mise" "activate" "pwsh" }}

# bunのglobalインストール先をPATHに追加
$bunBinPath = "{{ joinPath .chezmoi.homeDir ".bun" "bin" }}"
if ($env:PATH -notlike "*$bunBinPath*") {
    $env:PATH += ";$bunBinPath"
}

{{ end -}}

