# private_dot_config/powershell/profile.d/mise.ps1.tmpl
{{ if lookPath "mise" -}}
{{ output "mise" "activate" "pwsh" }}

# bunのglobalインストール先をPATHに追加
if ($env:PATH -notlike "*$($env:HOME)/.bun/bin*") {
    $env:PATH += ";$env:HOME/.bun/bin"
}

{{ end -}}

