# private_dot_config/powershell/profile.d/mise.ps1.tmpl
{{ if lookPath "mise" -}}
{{ output "mise" "activate" "pwsh" }}
{{ end -}}
