# private_dot_config/powershell/profile.d/oh-my-posh.ps1.tmpl
{{ output "oh-my-posh" "init" "pwsh" "--config=~/.config/powershell/theme.yaml" }}
