# private_dot_config/powershell/profile.d/vfox.ps1.tmpl
{{ if lookPath "vfox" -}}
{{ output "vfox" "activate" "pwsh" }}
{{ end -}}
