# private_dot_config/powershell/profile.d/zoxide.ps1.tmpl
{{ output "zoxide" "init" "--hook" "pwd" "powershell" }}

Set-Alias -Name cd -Value z -Option AllScope -Scope Global -Force
Set-Alias -Name cdi -Value zi -Option AllScope -Scope Global -Force
